# system/systemd/meuhypr-visor-boot-order.service
[Unit]
Description=MeuHypr — prioriza Visor Boot Manager no UEFI
After=local-fs.target
ConditionPathExists=/boot/efi/EFI/visor/visor_x64.efi

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/meuhypr-visor-boot-order
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
